# wsgi-0/myproject/myproject.service
[Unit]
Description=uWSGI server instance configured to server myproject

[Service]
Environment= PATH =/vagrant/myproject/myprojectenv/bin
WorkingDirectory=/vagrant/myproject
ExecStart=uwsgi --ini myproject.ini
